# detection-rules/local_rules.xml
<!--
  local_rules.xml - Custom Wazuh detection rules for the SOC lab.

  WHERE THIS GOES on the Wazuh manager:
      /var/ossec/etc/rules/local_rules.xml
  Then restart the manager:
      sudo systemctl restart wazuh-manager

  HOW WAZUH SEES THESE EVENTS:
      The Wazuh agent forwards Windows Sysmon events via the
      windows_eventchannel decoder. Sysmon fields become win.eventdata.*
      and Wazuh groups them as sysmon_event1 (process create),
      sysmon_event_13 (registry set), sysmon_event_10 (process access), etc.

  IMPORTANT - VERSION CHECK (honesty note):
      Parent group/rule names vary slightly between Wazuh versions. After you
      load these, run the matching Atomic Red Team test and confirm YOUR rule
      ID fires in the dashboard. If it does not, check Management > Rules for
      the correct Sysmon parent group on your version and adjust if_group.

  Custom rule IDs must be in the 100000-119999 range (reserved for user rules).
-->

<group name="local,sysmon,soc_lab,">

  <!-- ===========================================================
       T1059.001  Command and Scripting Interpreter: PowerShell
       Catches encoded / hidden PowerShell (a classic execution TTP).
       Data source: Sysmon Event ID 1 (Process Create).
       =========================================================== -->
  <rule id="100210" level="12">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image">powershell\.exe</field>
    <field name="win.eventdata.commandLine">-enc|-encodedcommand|-w hidden|-nop</field>
    <description>Encoded or hidden PowerShell command executed on $(win.system.computer)</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
  </rule>

  <!-- Higher-confidence variant: Office app spawning PowerShell (maldoc chain) -->
  <rule id="100213" level="14">
    <if_sid>100210</if_sid>
    <field name="win.eventdata.parentImage">WINWORD\.EXE|EXCEL\.EXE|POWERPNT\.EXE|OUTLOOK\.EXE</field>
    <description>Office application spawned encoded PowerShell - likely malicious document</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
  </rule>

  <!-- ===========================================================
       T1547.001  Registry Run Key persistence
       Catches new autostart entries under the Run keys.
       Data source: Sysmon Event ID 13 (Registry Value Set).
       =========================================================== -->
  <rule id="100211" level="10">
    <if_group>sysmon_event_13</if_group>
    <field name="win.eventdata.targetObject">\\CurrentVersion\\Run|\\CurrentVersion\\RunOnce</field>
    <description>Registry Run key persistence created on $(win.system.computer)</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <!-- ===========================================================
       T1003.001  OS Credential Dumping: LSASS Memory
       Catches a process opening a handle to lsass.exe (credential theft).
       Data source: Sysmon Event ID 10 (Process Access).
       =========================================================== -->
  <rule id="100212" level="13">
    <if_group>sysmon_event_10</if_group>
    <field name="win.eventdata.targetImage">lsass\.exe</field>
    <description>Suspicious access to LSASS memory on $(win.system.computer) - possible credential dumping</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
  </rule>

</group>
